# docker/supervisor.d/nginx.conf
[program:nginx]
command=/usr/sbin/nginx
process_name=%(program_name)s
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
